# crontab
0 0 * * * php artisan agenda:create 4
# An empty line is required at the end of this file for a valid cron file.
